# src/service/ros2_v2_mqtt_ros.service
[Unit]
Description="ROS2 node: mqtt_ros, version 2.0"
PartOf=ros2_v2_manage.service
After=ros2_v2_manage.service

[Service]
User=ubuntu
Type=simple
WorkingDirectory=/home/ubuntu/am_v2_ws
ExecStart=/bin/bash -c "source /opt/ros/galactic/setup.bash; source install/local_setup.bash; ros2 launch mqtt_ros mqtt_ros.launch.py"
KillSignal=SIGINT
Restart=on-failure
RestartSec=30

[Install]
WantedBy=ros2_v2_manage.service
